# build_deb/DEBIAN/postinst
#!/bin/bash
set -e
echo "Configuring Clack udev rules for ydotool..."
echo 'KERNEL=="uinput", GROUP="input", MODE="0660", OPTIONS+="static_node=uinput"' > /etc/udev/rules.d/99-uinput.rules
udevadm control --reload-rules
udevadm trigger

# Try to enable ydotoold service if it exists
if systemctl list-unit-files | grep -q ydotoold; then
    systemctl enable --now ydotoold || true
fi

if [ -n "$SUDO_USER" ]; then
    usermod -aG input "$SUDO_USER"
    echo "Added $SUDO_USER to the 'input' group. You may need to log out and back in."
fi

exit 0
